# server/dashboard-tools.service
[Unit]
Description=Dashboard Tools API (yt-dlp, rembg, PDF, file converter)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/dashboard
ExecStart=/usr/bin/python3 /opt/dashboard/tools_api.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Make sure /usr/local/bin is on PATH so subprocesses (deno, ffmpeg)
# get located when launched as www-data via systemd.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# rembg downloads its ONNX model on first use into ~/.u2net.
# www-data has /var/www as $HOME by default.
Environment=U2NET_HOME=/var/cache/dashboard-tools/u2net
# yt-dlp writes a small cache for some sites
Environment=XDG_CACHE_HOME=/var/cache/dashboard-tools
# Deno cache (used by the bgutil-ytdlp-pot-provider script-deno mode
# which generates YouTube PO Tokens — required for YT downloads on
# datacenter IPs to bypass the anti-bot challenge).
Environment=DENO_DIR=/var/cache/dashboard-tools/deno
# Playwright Chromium location — installed system-wide here so the
# www-data service user can launch the headless browser for the
# Reader Mode JS-rendering fallback.
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/playwright-browsers

[Install]
WantedBy=multi-user.target
